# neodct/overlay/etc/inittab
# /etc/inittab

# 1. Startup Logic
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -o remount,rw /
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mount -a
::sysinit:/sbin/swapon -a
::sysinit:/etc/init.d/rcS

# 2. The "Silent" GUI (Framebuffer)
# Run your script on the screen. No login prompt.
tty1::once:/bin/run_neodct.sh

# 3. The "Debug" Shell (Serial Port)
# getty_guard.sh spawns getty at 115200 only if the tty exists
# (ttyFIQ0 = Luckfox FIQ debugger, ttyAMA0 = QEMU PL011).
# Luckfox-flavor builds replace this file with inittab.luckfox (1500000).
::respawn:/bin/getty_guard.sh ttyFIQ0
::respawn:/bin/getty_guard.sh ttyAMA0


# 4. System Hooks
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
